# github-sync/crontab
MAILTO=hi@frozenridge.co
00 * * * * node ~/github-sync/index.js
