# github-sync/crontab
MAILTO=hi@frozenridge.co
00 * * * * (cd ~/github-sync && node index.js)
