# github-sync/crontab
MAILTO=hi@frozenridge.co
00 * * * * (cd ~/github-sync && /opt/node/bin/node index.js)
